# harness/launchd/com.solar.tech-hotspot-radar.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.solar.tech-hotspot-radar.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lisihao/Solar/harness/scripts/run_tech_hotspot_radar.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>UserName</key>
  <string>lisihao</string>

  <key>GroupName</key>
  <string>staff</string>

  <key>WorkingDirectory</key>
  <string>/Users/lisihao/Solar/harness</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONIOENCODING</key>
    <string>utf-8</string>
    <key>ALLOW_PARTIAL_SUCCESS</key>
    <string>1</string>
    <key>RUN_SOCIAL_REPORT</key>
    <string>0</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/lisihao/Solar/harness/logs/tech-hotspot-radar.daemon.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lisihao/Solar/harness/logs/tech-hotspot-radar.daemon.err.log</string>
</dict>
</plist>
