# solar/core/evolver/com.solar.task-executor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.solar.task-executor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sihaoli/Solar/core/evolver/task-executor.sh</string>
    </array>

    <!-- 每 5 分钟检查一次任务队列 -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/sihaoli/.solar/logs/task-executor.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sihaoli/.solar/logs/task-executor.stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
