# harness/com.solar.ai-influence-daily-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.solar.ai-influence-daily-digest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lisihao/Solar/harness/scripts/run_ai_influence_digest.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>8</integer>
      <key>Minute</key>
      <integer>7</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>13</integer>
      <key>Minute</key>
      <integer>7</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>18</integer>
      <key>Minute</key>
      <integer>7</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/lisihao/Solar/harness</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONIOENCODING</key>
    <string>utf-8</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/lisihao/Solar/harness/logs/ai-influence-daily-digest.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lisihao/Solar/harness/logs/ai-influence-daily-digest.err.log</string>
</dict>
</plist>
