# harness/com.solar.github-trends-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.solar.github-trends-digest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/lisihao/Solar/harness/scripts/run_github_trends_digest.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>37</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/lisihao/Solar/harness</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONIOENCODING</key>
    <string>utf-8</string>
    <key>GITHUB_TRENDS_MAIL_TO</key>
    <string>sean.lisihao@huawei.com</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/lisihao/Solar/harness/logs/github-trends-digest.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lisihao/Solar/harness/logs/github-trends-digest.err.log</string>
</dict>
</plist>
