# deploy/launchd/com.solar.security.system.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.solar.security.system</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SOLAR_HOME__/core/security/checks/system-check.sh</string>
    </array>

    <!-- 每 1 小时执行一次 -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__HOME_DIR__/.solar/logs/launchd-system.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME_DIR__/.solar/logs/launchd-system-error.log</string>
</dict>
</plist>
