# deploy/launchd/com.solar.security.access.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.solar.security.access</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sihaoli/Solar/core/security/checks/access-check.sh</string>
    </array>

    <!-- 每 1 分钟执行一次 (高频) -->
    <key>StartInterval</key>
    <integer>60</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/sihaoli/.solar/logs/launchd-access.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sihaoli/.solar/logs/launchd-access-error.log</string>
</dict>
</plist>
