# deploy/launchd/com.solar.security.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.solar.security.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sihaoli/Solar/core/security/checks/daily-audit.sh</string>
    </array>

    <!-- 每天凌晨 3 点执行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sihaoli/.solar/logs/launchd-daily.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/sihaoli/.solar/logs/launchd-daily-error.log</string>
</dict>
</plist>
