# deploy/com.solar.security.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.solar.security</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME_DIR__/.bun/bin/bun</string>
        <string>run</string>
        <string>__SOLAR_HOME__/core/security/security-daemon.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__SOLAR_HOME__/core/security</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME_DIR__</string>
        <key>PATH</key>
        <string>__HOME_DIR__/.bun/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME_DIR__/.solar/logs/security.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME_DIR__/.solar/logs/security-error.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
