# core/evolver/com.solar.native-evolver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.solar.native-evolver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__SOLAR_HOME__/core/evolver/native-evolver.sh</string>
    </array>

    <!-- 每天凌晨 3 点执行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__HOME_DIR__/.solar/logs/native-evolver.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME_DIR__/.solar/logs/native-evolver.stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
